# systemd/homeassistant-port-forward.service
[Unit]
Description=Private Home Assistant reverse tunnel and LM Studio/MQTT forward tunnels
After=network-online.target docker.service
Wants=network-online.target docker.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=encoded
EnvironmentFile=/etc/default/encodedpi-port-forward
ExecStart=/usr/bin/ssh -N \
    -p 9022 \
    -o ExitOnForwardFailure=yes \
    -o BatchMode=yes \
    -o ConnectTimeout=10 \
    -o ServerAliveInterval=30 \
    -o ServerAliveCountMax=3 \
    -R 127.0.0.1:8123:localhost:8123 \
    -R 127.0.0.1:8090:localhost:8090 \
    -R 127.0.0.1:8022:localhost:22 \
    -L 11434:127.0.0.1:1234 \
    -L 172.18.0.1:1884:127.0.0.1:1883 \
    ${REMOTE_TARGET}
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
